# zeropythia.service
[Unit]
Description=ZeroPythia – Three-phase zero-feed-in battery controller
Documentation=https://github.com/cnadler86/ZeroPythia
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zeropythia
Group=pythia
WorkingDirectory={{WORKING_DIR}}
ExecStart={{PYTHON_BIN}} main.py {{EXTRA_ARGS}}

Restart=always
RestartSec=10
TimeoutStopSec=30

# Logging via systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zeropythia

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths={{WORKING_DIR}} /var/cache/zeropythia-uv
ProtectHome=true
Environment=UV_CACHE_DIR=/var/cache/zeropythia-uv

[Install]
WantedBy=multi-user.target
